# supervisor.conf
[program:heatstop]
command=/usr/local/bin/python3.7 heatstop.py 
directory=/home/pi/heatstop
redirect_stderr=true
stopsignal=INT
environment=PYTHONUNBUFFERED="1"
